# roles/openjdk/templates/etc/profile.d/java.sh.j2
# {{ ansible_managed }}
export JAVA_HOME="{{ openjdk_home }}"
